# Systemd service/fancontrol.timer
[Unit]
Description=test

[Timer]
OnCalendar=*-*-* 06..18:*:*
Unit=fancontrol.service

[Install]
WantedBy=timers.target
